# docker/mariadb/conf/supervisor/conf.d/mysql.conf
[program:mysql]
command=/usr/bin/mysqld_safe
stdout_logfile=/var/log/supervisor/mysqld.log
autostart=true
autorestart=true
