# docker/php-fpm/conf/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.0 -F
stdout_logfile=/var/log/php-fpm.log
autostart=true
autorestart=true
