# docker/php-fpm/conf/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/sshd.log
autostart=true
autorestart=true
